# addons/neohub_control/rootfs/etc/services.d/neohub-control/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

cd /opt/neohub-control
exec /opt/neohub-control/venv/bin/python -m ha_addon.service
